# posos-github-repo/systemd/posos.service
[Unit]
Description=POSOS cash register
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=posos
Group=posos
Environment=DISPLAY=:0
Environment=POSOS_DATA_DIR=/var/lib/posos
ExecStart=/opt/posos/current/venv/bin/python -m posos
Restart=always
RestartSec=2

[Install]
WantedBy=graphical.target
